# systemd/jsonlogalert.service
[Unit]
Description=jsonlogalert.service

[Service]
Type=simple
ExecStart=/usr/local/bin/jsonlogalert
